# infra/ops-mcp/ops-mcp.service
[Unit]
Description=ops-mcp (scoped host operations for Atlas/Cowork)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ops-mcp
Group=ops-mcp
WorkingDirectory=/opt/ops-mcp
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/ops-mcp/env
ExecStart=/opt/ops-mcp/venv/bin/python /opt/ops-mcp/server.py
Restart=on-failure
RestartSec=5

# Hardening.
# NoNewPrivileges=no is required: the service runs sudo -n /bin/systemctl
# and sudo -n /bin/journalctl; those are SUID binaries that need to
# transition to UID 0. With NoNewPrivileges=yes the SUID bit is ignored.
# All other hardening directives are applied.
NoNewPrivileges=no
ProtectSystem=strict
# NOTE: PrivateTmp deliberately NOT set — per LESSONS #6, PrivateTmp causes
# cross-session visibility problems in similar MCP services. The service
# does not write to /tmp; its only writes are the audit log under /var/log.
ReadWritePaths=/var/log/ops-mcp /var/lib/ops-mcp
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
# Empty bounding set. sudo acquires UID 0 via the SUID bit (not via
# capability promotion), so an empty bounding set does not prevent the
# sudo→systemctl/journalctl hand-off from working. The subprocess runs
# as root UID but inherits no capabilities — sufficient for D-Bus systemctl
# and journal reads.
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
