# deploy/90-ventd-hwmon.rules
# ventd — grant the ventd group write access to hwmon pwm control files.
#
# The daemon runs unprivileged (User=ventd, Group=ventd; see
# deploy/ventd.service). It relies on Discretionary Access Control to
# open pwm<N> and pwm<N>_enable for write; this rule arranges that DAC
# by chgrp'ing those files to the ventd group with g+w whenever a
# hwmon device appears — on boot, hot-plug, driver rebind, or after
# an `rmmod; modprobe` cycle.
#
# Installed to:
#   /etc/udev/rules.d/90-ventd-hwmon.rules
#
# ──────────────────────────────────────────────────────────────────────
#
# Design: chip-agnostic.
#
# This rule does not key on ATTR{name}=="<chip>". It fires on every
# hwmon device and grants group access to whatever pwm<N> /
# pwm<N>_enable files the chip exposes, if any. Temperature-only chips
# (coretemp, k10temp, drivetemp, acpitz) are traversed but nothing
# happens because the `pwm[0-9]*` glob expands to no matches. The
# chgrp and chmod commands each have `2>/dev/null` to suppress the
# "No such file" error those no-match cases produce.
#
# This is deliberate and non-negotiable. ventd's install promise is
# "zero terminal after install" — keying on chip name would require
# the operator to identify their chip and edit this file. The
# alternative (shipping a pre-populated list of chip names) is a
# maintenance trap: every new board generation brings a new Super-I/O
# variant, and the list goes stale.
#
# ──────────────────────────────────────────────────────────────────────
#
# Safety envelope
#
# - Group write only (g+w). Never world-writable, never setuid.
# - Subsystem filter (SUBSYSTEM=="hwmon") ensures this rule only
#   touches hwmon devices — not block devices, input devices, USB,
#   or anything else udev might see.
# - Scope limited to pwm<N> and pwm<N>_enable. Temperature inputs,
#   fan RPM inputs, voltage inputs, power inputs are untouched; their
#   default read-for-all, write-for-root semantics are preserved.
# - The rule grants the same DAC envelope whether the hwmon device
#   is a motherboard Super-I/O chip, an AIO pump controller, an
#   amdgpu card, or a future driver that has not shipped yet.
#
# ──────────────────────────────────────────────────────────────────────
#
# Reload after editing:
#   sudo udevadm control --reload
#   sudo udevadm trigger --subsystem-match=hwmon
#
# The installer (scripts/install.sh and scripts/postinstall.sh) runs
# both commands automatically after placing this file, so no operator
# action is required on a clean install.

SUBSYSTEM=="hwmon", RUN+="/bin/sh -c 'chgrp ventd /sys%p/pwm[0-9]* /sys%p/pwm[0-9]*_enable 2>/dev/null; chmod g+w /sys%p/pwm[0-9]* /sys%p/pwm[0-9]*_enable 2>/dev/null; exit 0'"
