# deploy/github-mcp.service
[Unit]
Description=GitHub MCP Server (HTTP mode, self-hosted, ci-log-tools build)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
ExecStartPre=-/usr/bin/docker rm -f github-mcp
ExecStart=/usr/bin/docker run --rm --name github-mcp \
    --env-file /etc/github-mcp/env \
    -p 127.0.0.1:8765:8082 \
    github-mcp-server:ci-log-tools \
    http
ExecStop=/usr/bin/docker stop github-mcp

[Install]
WantedBy=multi-user.target
