# scripts/ventd.service
[Unit]
Description=ventd — automatic Linux fan control daemon
Documentation=https://github.com/ventd/ventd
# Modules in /etc/modules-load.d/ are loaded by systemd-modules-load.
# ventd writes its winning module there on first boot, so subsequent
# boots bring up hwmon before this service starts.
After=network.target systemd-modules-load.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ventd --config /etc/ventd/config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Restart on any failure — covers crashes, OOM kills, and bad exits.
Restart=always
# Wait 3 s before restarting so a rapid crash loop doesn't hammer the hardware.
RestartSec=3s

# Run as root: required to write /sys/class/hwmon/*/pwm* and call modprobe.
User=root

# Capture stdout/stderr in journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ventd

[Install]
WantedBy=multi-user.target
