# .cowork/tools/spawn-mcp/spawn-mcp.service
[Unit]
Description=spawn-mcp (Cowork -> Claude Code dispatcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cc-runner
Group=cc-runner
WorkingDirectory=/opt/spawn-mcp
Environment=PYTHONUNBUFFERED=1
Environment=IS_DEMO=1
# spawn_cc_batch: off by default; enable once IMPROV-A is deployed and verified
Environment=SPAWN_MCP_BATCH_ENABLED=false
EnvironmentFile=/etc/spawn-mcp/env
ExecStart=/opt/spawn-mcp/venv/bin/python /opt/spawn-mcp/server.py
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
# ProtectHome deliberately NOT set. The service runs as cc-runner and
# `claude` needs to read+write /home/cc-runner/.claude/ for session state,
# cache, and auth tokens. Since the service IS cc-runner, there is no
# boundary to enforce; ProtectHome=read-only here would only block the
# service's own user from writing its own home.
ReadWritePaths=/var/log/spawn-mcp /tmp
# NOTE: PrivateTmp deliberately NOT set. Operators attach to live sessions
# with `sudo -u cc-runner tmux attach -t cc-...` from an external shell,
# and tail_session needs to reach the tmux server across spawn-mcp
# restarts. PrivateTmp would put tmux's socket and our prompt files in a
# namespace-private tmpfs that unmounts on restart, breaking both paths.
# Prompt files are written 0600 into /tmp/spawn-mcp/ and include a SHA
# returned to Cowork so tampering is detectable.
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
# NOTE: MemoryDenyWriteExecute deliberately NOT set. spawn-mcp itself is
# pure Python and would comply, but several recent CPython builds JIT
# regex fragments via PEP 659 specializations that touch W+X pages. If
# this service ever trips SIGSEGV on startup, this is the first knob to
# investigate.
RestrictRealtime=yes
SystemCallArchitectures=native

# No cross-user handoff, so no capabilities required.
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
