# .cowork/tools/spawn-mcp/spawn-mcp.service
[Unit]
Description=spawn-mcp (Cowork -> Claude Code dispatcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spawn-mcp
Group=spawn-mcp
WorkingDirectory=/opt/spawn-mcp
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/spawn-mcp/env
ExecStart=/opt/spawn-mcp/venv/bin/python /opt/spawn-mcp/server.py
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/spawn-mcp /tmp/cc-runner
# NOTE: PrivateTmp deliberately NOT set. spawn-mcp writes prompt files to
# /tmp/cc-runner/<session>.md and cc-runner's tmux session reads them —
# this cross-user handoff requires both processes to see the same /tmp.
# Prompts are chmod 600 immediately after write, owned by cc-runner,
# and the directory itself is chmod 700 owned by cc-runner. A prompt SHA
# is returned to Cowork so tampering is detectable.
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
# NOTE: MemoryDenyWriteExecute deliberately NOT set. spawn-mcp itself is
# pure Python and would comply, but several recent CPython builds JIT
# regex fragments via PEP 659 specializations that touch W+X pages. If
# this service ever trips SIGSEGV on startup, this is the first knob to
# investigate.
RestrictRealtime=yes
SystemCallArchitectures=native

# spawn-mcp needs to exec `sudo -u cc-runner tmux ...`
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
