# .cowork/tools/spawn-mcp/spawn-mcp-tunnel.service
[Unit]
Description=cloudflared tunnel for spawn-mcp
# Wants+After — tunnel should start with spawn-mcp, but a restart of
# spawn-mcp.service MUST NOT cascade-restart the tunnel, because that
# rotates the trycloudflare quick-tunnel hostname and breaks the
# claude.ai connector wiring. Requires= would cascade.
After=spawn-mcp.service
Wants=spawn-mcp.service

[Service]
Type=simple
User=spawn-mcp
ExecStart=/usr/local/bin/cloudflared tunnel --url http://127.0.0.1:8891 --no-autoupdate
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
